# hipro-api.service
[Unit]
Description=HiPro Commerce API Server
Documentation=https://github.com/Amit-Kumar22/hipro-comm
After=network.target mongod.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/hipro-comm-api
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hipro-api
Environment=NODE_ENV=production
Environment=PORT=5000
Environment=MONGODB_URI=mongodb://localhost:27017/hipro-comm-prod
Environment=JWT_SECRET=your-production-jwt-secret-change-this
Environment=JWT_EXPIRES_IN=7d
Environment=COOKIE_SECRET=your-production-cookie-secret-change-this
Environment=FRONTEND_URL=https://shop.hiprotech.org
Environment=ADMIN_URL=https://adminshop.hiprotech.org
Environment=RATE_LIMIT_WINDOW_MS=900000
Environment=RATE_LIMIT_MAX_REQUESTS=100
Environment=UPLOAD_PATH=uploads
Environment=MAX_FILE_SIZE=10485760
Environment=ALLOWED_FILE_TYPES=image/jpeg,image/png,image/webp,image/gif
Environment=LOG_LEVEL=info
Environment=EMAIL_SERVICE=smtp
Environment=EMAIL_USER=info@hiprotech.org
Environment=EMAIL_PASS=Abhi@2026
Environment=EMAIL_FROM=info@hiprotech.org
Environment=SMTP_HOST=smtp.hostinger.com
Environment=SMTP_PORT=587
Environment=SMTP_SECURE=false
Environment=SMTP_USER=info@hiprotech.org
Environment=SMTP_PASS=Abhi@2026
Environment=SMTP_FROM_NAME=HiPro Commerce
Environment=SMTP_FROM_EMAIL=info@hiprotech.org

[Install]
WantedBy=multi-user.target